# portal_login.service
[Unit]
Description=Auto WiFi Login for College Captive Portal
After=network.target

[Service]
ExecStart=/usr/bin/python3 /home/username/wifi_login.py
Restart=on-failure
User=your_username
WorkingDirectory=/home/username
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
